# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:backend]
; Try running uvicorn directly without shell wrapper
command=uvicorn backend.main:app --host 0.0.0.0 --port 8000 --log-level debug
directory=/app
environment=PYTHONPATH="/app",PATH="/usr/local/bin:/usr/bin:/bin"
stdout_logfile=/var/log/supervisor/backend.log
stderr_logfile=/var/log/supervisor/backend.err
autorestart=true
startsecs=5
startretries=10
user=root

[program:nginx]
command=nginx -g 'daemon off;'
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx.err
autorestart=true
startsecs=3
